# testing/display/display_config.service
[Unit]
Description=Display Configuration Service
After=network.target

[Service]
User=your_username
WorkingDirectory=/path/to/script/directory
ExecStart=/usr/bin/python3 /path/to/script/display_config.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

